# scripts/templates/v3.2.x/support-files/templates/#etc#supervisor-cmdb-server.conf
[unix_http_server]
file=__BK_HOME__/logs/cmdb/supervisor-cmdb3.sock


[supervisorctl]
serverurl=unix://__BK_HOME__/logs/cmdb/supervisor-cmdb3.sock ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=__BK_HOME__/logs/cmdb/supervisord_cmdb3.log          ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB                                           ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10                                              ; (num of main logfile rotation backups;default 10)
loglevel=info                                                   ; (log level;default info; others: debug,warn,trace)
pidfile=__BK_HOME__/logs/cmdb/supervisord_cmdb.pid          ; (supervisord pidfile;default supervisord.pid)
nodaemon=false                                                  ; (start in foreground if true;default false)
minfds=1024                                                     ; (min. avail startup file descriptors;default 1024)
minprocs=200                                                    ; (min. avail process descriptors;default 200)
;umask=022                                                      ; (process file creation umask;default 022)
;user=root                                                      ; (default is current user, required if root)
;identifier=supervisor                                          ; (supervisord identifier, default is 'supervisor')
;directory=/tmp/                                                ; (default is not to cd during start)
;nocleanup=true                                                 ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp/                                              ; ('AUTO' child log dir, default $TEMP)
;environment=KEY="value"                                        ; (key value pairs to add to environment)
;strip_ansi=false

[program:cmdb_auditcontoller]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_auditcontroller  --addrport=__LAN_IP__:__CMDB_AUDIT_CTRL_PORT__  --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_auditcontoller.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_auditcontoller.stderr.log


[program:cmdb_hostcontroller]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_hostcontroller  --addrport=__LAN_IP__:__CMDB_HOST_CTRL_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_hostcontroller.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_hostcontroller.stderr.log

[program:cmdb_objectcontroller]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_objectcontroller  --addrport=__LAN_IP__:__CMDB_OBJECT_CTRL_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_objectcontroller.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_objectcontroller.stderr.log

[program:cmdb_proccontroller]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_proccontroller  --addrport=__LAN_IP__:__CMDB_PROC_CTRL_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_proccontroller.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_proccontroller.stderr.log


[program:cmdb_eventserver]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_eventserver  --addrport=__LAN_IP__:__CMDB_EVENT_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3  --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_eventserver.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_eventserver.stderr.log

[program:cmdb_hostserver]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_hostserver  --addrport=__LAN_IP__:__CMDB_HOST_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_hostserver.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_hostserver.stderr.log



[program:cmdb_adminserver]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_adminserver  --addrport=__LAN_IP__:__CMDB_ADMIN_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --config=__BK_HOME__/cmdb/server/conf/migrate.conf 
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_adminserver.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_adminserver.stderr.log


[program:cmdb_procserver]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_procserver  --addrport=__LAN_IP__:__CMDB_PROC_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_procserver.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_procserver.stderr.log



[program:cmdb_toposerver]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_toposerver  --addrport=__LAN_IP__:__CMDB_TOPO_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_toposerver.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_toposerver.stderr.log

[program:cmdb_apiserver]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_apiserver  --addrport=__LAN_IP__:__CMDB_API_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_apiserver.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_apiserver.stderr.log


[program:cmdb_datacollection]
username=cmdb
directory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_datacollection  --addrport=__LAN_IP__:__CMDB_DATACOLLECTION_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3 --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_datacollection.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_datacollection.stderr.log

[program:cmdb_webserver]
username=cmdb
diectory=__BK_HOME__/cmdb/
command= __BK_HOME__/cmdb/server/bin/cmdb_webserver  --addrport=__LAN_IP__:__CMDB_WEB_PORT__ --logtostderr=false --log-dir=__BK_HOME__/logs/cmdb --v=3  --regdiscv=__ZK_HOST__:__ZK_PORT__
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/cmdb/cmdb_webserver.stdout.log
stderr_logfile=__BK_HOME__/logs/cmdb/cmdb_webserver.stderr.log
